# conf/supervisor/gateway.conf
[program:gateway]
command=/usr/sbin/nginx
autostart=true
autorestart=true
user=root
stdout_logfile=/dev/null
redirect_stderr=true
stdout_logfile_maxbytes=20
